# ddns-updater(automated)/automation (systemd)/dns_updater.service
[Unit]
Description=Cloudflare DDNS Updater
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=your-username
WorkingDirectory=/opt/ddns-updater
ExecStart=/opt/ddns-updater/venv/bin/python /opt/ddns-updater/main.py
EnvironmentFile=/etc/ddns-updater/.env
Environment="PYTHONUNBUFFERED=1"
Restart=on-failure
RestartSec=10
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target